# scripts/debian/xgrid.cron.d
0/5 * * * * root xgrid-graph >/dev/null 2>&1
#
